# etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:couchdb]
command=/usr/bin/couchdb

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
